# systemd/hltvf-scrap.timer
[Unit]
Requires=hltvf-scrap.service

[Timer]
OnCalendar=*-*-* 4..21:55:00
RandomizedDelaySec=180

[Install]
WantedBy=timers.target
