# crond/crontab
0 5 * * * /???/start-notifier.sh
55 4-21 * * * /???/start-scraper.sh
